# conf/supervisor/nginx.ini
[program:nginx]
command=/usr/local/nginx/sbin/nginx -c /usr/local/nginx/conf/nginx.conf
process_name=%(program_name)s
autostart=true
autorestart=true
user=root
startsecs=0
redirect_stderr=true
stdout_logfile=/opt/logs/supervisor_nginx.log
stdout_logfile_maxbytes=0
