# conf/supervisor/php-fpm.ini
[program:php-fpm]
command=/usr/local/php/sbin/php-fpm
process_name=%(program_name)s_%(process_num)
autostart=true
autorestart=true
user=www
startsecs=0
redirect_stderr=true
stdout_logfile=/opt/logs/php-fpm.log
stdout_logfile_maxbytes=0
